# roles/base/templates/journal-vacuum.service
[Unit]
Description=Vacuum the journal logs
After=network.target

[Service]
ExecStart=journalctl --vacuum-size=80M

[Install]
WantedBy=default.target
